# systemd/flirt-deploy.timer
[Unit]
Description=Run Flirt pull-based deploy every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
